# .config/systemd/user/picom.path
[Unit]
Description=watch for Picom config

[Path]
PathChanged=%E/picom/picom.conf
Unit=watcher@picom.service

[Install]
WantedBy=default.target
