# .config/systemd/user/redshift.path
[Unit]
Description=watch for Redshift config

[Path]
PathChanged=%E/redshift/redshift.conf
Unit=watcher@redshift.service

[Install]
WantedBy=default.target
